# scripts/com.stockcontest.weekly-email.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.stockcontest.weekly-email</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/Yitzi/Library/Application Support/stock-contest/scripts/weekly-email.sh</string>
    </array>

    <!--
      Drift-proof scheduling: poll every 5 minutes and let weekly-email.sh decide
      whether it is inside the send window (Friday >=16:45 .. Sun, once per ISO week).
      We deliberately do NOT use StartCalendarInterval: on this laptop, sleep-deferred
      calendar fires get re-anchored by macOS UserEventAgent and drift to the wrong
      time (they had crept to ~9:45 AM). StartInterval has no calendar anchor to drift;
      non-window polls exit in milliseconds via the script's guard.
    -->
    <key>StartInterval</key>
    <integer>300</integer>

    <key>RunAtLoad</key>
    <true/>

    <key>StandardOutPath</key>
    <string>/Users/Yitzi/Library/Application Support/stock-contest/logs/launchd-weekly-email-stdout.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/Yitzi/Library/Application Support/stock-contest/logs/launchd-weekly-email-stderr.log</string>

    <key>ProcessType</key>
    <string>Standard</string>
</dict>
</plist>
