# scripts/com.stockcontest.daily-refresh.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.stockcontest.daily-refresh</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/Yitzi/Library/Application Support/stock-contest/scripts/daily-refresh.sh</string>
    </array>

    <!--
      Drift-proof scheduling: poll every 5 minutes and let daily-refresh.sh decide
      whether it is inside a refresh slot (weekday AM >=09:31, PM >=16:20, once per
      slot per day). We deliberately do NOT use StartCalendarInterval: on this laptop,
      sleep-deferred calendar fires get re-anchored by macOS UserEventAgent and drift
      to the wrong time (the slots had drifted to ~02:31 AM / 09:20 AM). StartInterval
      has no calendar anchor to drift; out-of-window polls exit in milliseconds.
    -->
    <key>StartInterval</key>
    <integer>300</integer>

    <key>RunAtLoad</key>
    <true/>

    <key>StandardOutPath</key>
    <string>/Users/Yitzi/Library/Application Support/stock-contest/logs/launchd-daily-refresh-stdout.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/Yitzi/Library/Application Support/stock-contest/logs/launchd-daily-refresh-stderr.log</string>

    <key>ProcessType</key>
    <string>Standard</string>
</dict>
</plist>
